# profile.d/shell.sh
export PS1='[\u@\h#${CONTAINER_NAME} \W]\$ '
